# CDS-Auto-Announce/deploy/ubuntu22.04/ip-announce-web.service
[Unit]
Description=IP Segment Announce/Withdraw Web Console
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/IP段自动宣告-20260602
Environment=IP_ANNOUNCE_CONFIG=/home/ubuntu/IP段自动宣告-20260602/config.yaml
Environment=IP_ANNOUNCE_HOST=0.0.0.0
Environment=IP_ANNOUNCE_PORT=9010
ExecStart=/home/ubuntu/IP段自动宣告-20260602/.venv/bin/gunicorn -w 2 -b 0.0.0.0:9010 --timeout 900 --graceful-timeout 30 --chdir /home/ubuntu/IP段自动宣告-20260602 wsgi:app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
